# share/soft/jenkins/jenkins.service
[Unit]
Description=Jenkins Daemon

[Service]
User=jenkins
ExecStart=/usr/local/bin/jenkins.sh

[Install]
WantedBy=multi-user.target